# config/docker/app/supervisord.d/apiServer.ini
[program:apiServer]
command = /home/bitmex/api/console/runner run apiServer
directory = /home/bitmex/api
autorestart = true
startsecs = 10
user = bitmex
stopasgroup = true
stopsignal = SIGTERM
environment = LC_LANG="en_US.utf8",LC_ALL="en_US.utf8",pipenv_path=/home/bitmex/.local/bin/pipenv,HOME=/home/bitmex/